# deploy/hermes-datamanagementd.service
[Unit]
Description=Hermes Data Management Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=hermes
Group=hermes
WorkingDirectory=/opt/hermes
ExecStart=/opt/hermes/datamanagementd \
  -socket-path /tmp/hermes-datamanagement.sock \
  -sqlite-path /var/lib/hermes/datamanagement/datamanagementd.db \
  -version 1.0.0
Restart=always
RestartSec=5s
LimitNOFILE=100000
NoNewPrivileges=true
PrivateTmp=false

[Install]
WantedBy=multi-user.target
